# images/dev-toolchain/services.d/dockerd/run
#!/command/execlineb -P
#

/usr/bin/dockerd
